# udev/rules.d/92-dmabuf-peer-mem.rules
ACTION=="add", SUBSYSTEM=="infiniband", RUN+="/sbin/modprobe dmabuf_peer_mem"
KERNEL=="dmabuf_reg", MODE="0666"
